# extras/udev/rules.d/99-vibranium-ac.rules
SUBSYSTEM=="power_supply", ATTR{type}=="Mains", ACTION=="change", RUN+="/usr/local/bin/vibranium-core-ac-notification"

# vim:ft=udev
