# defaults/profile.d/20-amdgpu.sh
# Explicitly set propper GPU drivers
# https://wiki.archlinux.org/title/AMDGPU
export VDPAU_DRIVER="radeonsi"
export LIBVA_DRIVER_NAME="radeonsi"
