# defaults/profile.d/20-sudo-prompt.sh
export SUDO_PROMPT="$(tput setaf 1)$(tput bold)[sudo]$(tput sgr0) Password for $USER: "

